# device_service/card-wifi-watchdog.service
[Unit]
Description=Card Feeder — WiFi watchdog (auto AP when no home WiFi)
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
WorkingDirectory=/home/bbsw/workspace/device_service
ExecStart=/usr/bin/python3 -u /home/bbsw/workspace/device_service/wifi_watchdog.py
Environment=CARD_WIFI_IFACE=wlan0
Environment=CARD_WIFI_GRACE=40
# AP tự hạ sau 180s nếu không ai cấu hình (an toàn khi test từ xa).
# Sản xuất thực thụ: đặt CARD_AP_TIMEOUT=0 để AP giữ tới khi cấu hình xong.
Environment=CARD_AP_TIMEOUT=180
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
